# .docker/supervisord.conf
[program:uwsgi]
user = uwsgi
command=/opt/bin/uwsgi --ini 
autostart=true
autorestart=true
stderr_logfile = 
stdout_logfile = 
stopsignal=INT
